# dot_config/systemd/user/mako-enable-dnd.service
[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/makoctl mode -a do-not-disturb
