# dot_config/systemd/user/mako-disable-dnd.service
[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/makoctl mode -r do-not-disturb
